# releng/airootfs/etc/profile.d/launch-calamares.sh
#!/bin/bash
[ -z "$DISPLAY" ] || ( sleep 2 && calamares & )
